# demo_ansible/new_ansible_folder/files/dell_monitor.service
[Unit]
Description=...

[Service]
Type=oneshot
RemainAfterExit=true
ExecStop=/usr/local/bin/dell_monitor.sh

[Install]
WantedBy=multi-user.target
