# demo_ansible/new_ansible_folder/files/dell_monitor.service
[Unit]
Before=shutdown.target reboot.target
After=mysql.service
Requires=network.target

[Service]
ExecStart=/usr/local/bin/dell_monitor.sh

[Install]
WantedBy=default.target
